# private_dot_config/private_arch/ansible/roles/web/templates/openclaw-certbot-renew.timer.j2
[Unit]
Description=Twice-daily Let's Encrypt renewal timer for the OpenClaw reverse proxy

[Timer]
OnCalendar=*-*-* 03,15:00:00
RandomizedDelaySec=30m
Persistent=true

[Install]
WantedBy=timers.target
