# private_dot_config/private_arch/ansible/roles/services/templates/media-sorter-process.timer.j2
[Unit]
Description=Run media sorter queue processor

[Timer]
OnBootSec=2min
OnUnitActiveSec=1min
RandomizedDelaySec=15s
Unit=media-sorter-process.service

[Install]
WantedBy=timers.target
