# private_dot_config/private_arch/ansible/roles/services/templates/media-sorter-process.service.j2
[Unit]
Description=Process queued Transmission media sorter items
After=network-online.target transmission.service
Wants=network-online.target

[Service]
Type=oneshot
User=transmission
Group=transmission
UMask=0002
EnvironmentFile=-{{ arch_media_sorter_tmdb_env_file }}
EnvironmentFile=-{{ arch_media_sorter_telegram_env_file }}
EnvironmentFile=-{{ arch_media_sorter_xai_env_file }}
Environment=MEDIA_SORTER_TELEGRAM_ENV_FILE={{ arch_media_sorter_telegram_env_file }}
ExecStart={{ arch_media_sorter_path }} --process-queue --queue-root {{ arch_media_sorter_queue_root }} --series-root {{ arch_media_series_dir }} --films-root {{ arch_media_films_dir }} --music-root {{ arch_media_music_dir }} --source-root {{ arch_transmission_download_dir }} --audiodb-api-key {{ arch_media_sorter_audiodb_api_key }}{% if arch_media_sorter_grok_review_enabled | bool %} --grok-review --xai-responses-url {{ arch_media_sorter_xai_responses_url }} --xai-model {{ arch_media_sorter_xai_model }} --xai-timeout {{ arch_media_sorter_xai_timeout }}{% endif %}
